# services/mempalace/com.shishka.mempalace-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  MemPalace nightly backup — macOS launchd plist.
  Spec: docs/plans/spec-shishka-brain.md §3.3
  MC task: c2778157 (Brain: auto-ingest pipeline)

  Install:
    cp services/mempalace/com.shishka.mempalace-backup.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.shishka.mempalace-backup.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.shishka.mempalace-backup.plist
    rm ~/Library/LaunchAgents/com.shishka.mempalace-backup.plist

  Manual run:
    launchctl start com.shishka.mempalace-backup

  Logs:
    tail -f ~/Library/Logs/mempalace-backup.log
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shishka.mempalace-backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/lesianich/Library/CloudStorage/GoogleDrive-lesia@shishka.health/Общие диски/Shishka healthy kitchen/services/mempalace/backup.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/mempalace-backup.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/mempalace-backup.log</string>

  <key>RunAtLoad</key>
  <false/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
